# foundry/systemd/foundry-refresh.timer
[Unit]
Description=Daily Foundry ledger refresh

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
